# scripts_timers/package_db_refresh.timer
[Unit]
Description=Run date.service every 15 minutes

[Timer]
OnCalendar=*:0/15
